# files/etc/crontabs/root
# Crontab for OMP Schedule tasks
0 * * * *   omp-run-scheduled
